# debian/postinst
update-rc.d zrpc defaults
/etc/init.d/zrpc start
